# monay.service
[Unit]
Description=MonAY Finance Video System
After=network.target
Wants=network-online.target
StartLimitIntervalSec=0

[Service]
Type=simple
User=monay
Group=monay
WorkingDirectory=/opt/monay
Environment=PATH=/opt/monay/venv/bin:/usr/local/bin:/usr/bin:/bin
Environment=VIRTUAL_ENV=/opt/monay/venv
Environment=PYTHONPATH=/opt/monay
Environment=PYTHONUNBUFFERED=1
Environment=WAN_VENV_PATH=/opt/monay/wan_venv
Environment=MAIN_VENV_PATH=/opt/monay/venv
ExecStart=/opt/monay/venv/bin/python -m src.enhanced_main
Restart=always
RestartSec=30
StandardOutput=journal
StandardError=journal
SyslogIdentifier=monay
TimeoutStartSec=300
TimeoutStopSec=300

# Resource limits (CPU-only server optimized)
LimitNOFILE=65536
MemoryMax=4G
CPUQuota=50%

# Security settings
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/monay
PrivateTmp=true

[Install]
WantedBy=multi-user.target